# packaging/launchd/com.blin.arcmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  arcmux daemon LaunchAgent (macOS).

  This is a template: @HOME@ / @BIN@ / @LOG_FILE@ / @PATH@ are substituted by
  `make service-install`, which renders the result to
  ~/Library/LaunchAgents/com.blin.arcmux.plist and bootstraps it into the
  per-user launchd domain. See packaging/launchd/README.md.

  RunAtLoad + KeepAlive make launchd the supervisor (auto-start at login/boot,
  restart on crash) — the macOS analogue of the systemd Restart=always unit on
  labs. The Makefile detects this agent (HAS_LAUNCHD_AGENT) and routes
  start/stop/restart through launchctl so it never fights KeepAlive.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.blin.arcmux</string>

  <key>ProgramArguments</key>
  <array>
    <string>@BIN@</string>
    <string>start</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>@PATH@</string>
    <key>TZ</key>
    <string>America/Los_Angeles</string>
  </dict>

  <key>StandardOutPath</key>
  <string>@LOG_FILE@</string>
  <key>StandardErrorPath</key>
  <string>@LOG_FILE@</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
